# src/GPIO_interface/scripts/gpio
#!/bin/sh /etc/rc.common
# GPIO startup script

START=1
STOP=2

start() {
        /intelli/gpio/gpio_interface /intelli/gpio/config/gpio_config.xml &
}

stop() {
        killall gpio_interface
}

boot() {
        /intelli/gpio/gpio_interface /intelli/gpio/config/gpio_config.xml &
}
